# conf/varnish/supervisord.conf
[supervisord]
nodaemon=true

[program:varnish]
command = /usr/sbin/varnishd -f /etc/varnish/default.vcl
user = root
autostart = true
