# conf/gearman/supervisord.conf
[supervisord]
nodaemon=true

[program:gearman]
command = /usr/sbin/gearmand --user=gearman --log-file=/var/log/gearman-job-server/gearman.log  
user = root
autostart = true
